# devices/x86_64/factory/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# update chnroute and gfwlist
if [ ! -d /root/update_chnroute_dnsmasqconf_list.sh ]; then
    sh /root/update_chnroute_dnsmasqconf_list.sh > /dev/null 2>&1
fi

exit 0
